# assets/packages/hareply.service
[Unit]
Description=hareply
After=syslog.target network-online.target

[Service]
EnvironmentFile=/etc/hareply/hareply.conf
ExecStart=/usr/bin/hareply serve
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target